# templates/systemd/syncmail.service.j2
[Unit]
Description=Synchronize mail dirs between IMAP accounts and this computer

[Service]
ExecStart=/usr/bin/offlineimap -o -l /tmp/offlineimap.log
Type=oneshot
