# templates/systemd/syncmail.timer.j2
[Unit]
Description=Synchronize mail every 10 minutes and on reboot

[Timer]
OnBootSec=20s
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
